# systemd/homelab-backup-cloud.service
[Unit]
Description=Homelab Cloud Backup Service (Yandex Disk)
Documentation=https://github.com/vaultwarden/server
After=network.target docker.service
Requires=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/srv/homelab/homelab-server
ExecStart=/srv/homelab/homelab-server/scripts/backup.sh cloud
Nice=15
IOSchedulingClass=idle
IOSchedulingPriority=7

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
